# cookbook/templates/default/couch_watcher/who_watches_the_couch_watcher.service.erb
[Unit]
Description=Who Watches The Couch Watcher Restart Service
After=network.target

[Service]
Type=simple
Restart=always
PrivateTmp=yes
User=root
WorkingDirectory=<%= node[:primero][:app_dir] %>
ExecStart=<%= node[:primero][:daemons_dir] %>/who-watches-the-couch-watcher.sh

[Install]
WantedBy=multi-user.target
