# cookbook/files/default/xvfb/upstart_script.conf
# Xvfb display for browser testing

description     "Xvfb service"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec Xvfb :10 -ac > /dev/null
