# docker/db/supervisord.conf
[supervisord]

[program:couchdb]
command = /usr/bin/couchdb
autostart = true
user = couchdb

[program:nginx]
command = /usr/sbin/nginx
autorestart = true
autostart = true
