# cookbook/templates/default/passenger/passenger.service.erb
[Unit]
Description=Passenger Standalone Application Server
After=network.target

[Service]
Type=forking
Restart=always
PrivateTmp=yes
User=<%= node[:primero][:app_user] %>
Group=<%= node[:primero][:app_group] %>
WorkingDirectory=<%= node[:primero][:app_dir] %>
PIDFile=<%=@passenger_pid %>
ExecStart=<%= @passenger_app %> start
ExecReload=<%= @passenger_app %> reload
ExecStop=<%= @passenger_app %> stop

[Install]
WantedBy=multi-user.target
