# docker/nginx/root/etc/cron.d/certbot-renew
0 9 * * * certbot renew --post-hook "kill -s HUP $(cat /var/run/nginx.pd)"
