# cookbook/templates/default/puma/puma.service.erb
[Unit]
Description=Puma Application Server
After=network.target

[Service]
Type=simple
Restart=always
PrivateTmp=yes
User=<%= node[:primero][:app_user] %>
Group=<%= node[:primero][:app_group] %>
WorkingDirectory=<%= node[:primero][:app_dir] %>
PIDFile=<%= @puma_pid %>
Environment="RAILS_ENV=<%= node[:primero][:rails_env] %>"
Environment="RAILS_LOG_PATH=<%= @rails_log_dir %>"
export CHECK_ONLINE_STATUS=<%= node[:primero][:check_online_status] %>
ExecStart=<%= node[:primero][:home_dir] %>/.rvm/bin/rvm default do bundle exec puma -C <%= node[:primero][:app_dir] %>/config/puma.rb
ExecStop=<%= node[:primero][:home_dir] %>/.rvm/bin/rvm default do bundle exec pumactl -S <%= @puma_state %> stop

[Install]
WantedBy=multi-user.target
