# cookbook/templates/default/data-mounts.conf.erb
description     "Mount Primero data filesystems on boot"
start on local-filesystems
console log

script
  DATA_PART=$(blkid -L primero-data)

  if [ ! -z "$DATA_PART" ]
  then
    echo 'primero-data filesystem found'
    mkdir -p /mnt/primero-data
    mount $DATA_PART /mnt/primero-data
    echo "$DATA_PART mounted at /mnt/primero-data"

    # Make blank dirs to move to simplify logic
    mkdir -p <%= node[:primero][:log_dir] %> /var/lib/couchdb
    echo "Moving existing dirs if none exist on data filesystem"
    [ -d /mnt/primero-data/logs ] || mv <%= node[:primero][:log_dir] %> /mnt/primero-data/logs
    [ -d /mnt/primero-data/couchdb ] || mv /var/lib/couchdb /mnt/primero-data/couchdb

    # Remake the empty dirs to provide a mount target
    mkdir -p <%= node[:primero][:log_dir] %> /var/lib/couchdb
    echo "Bind mounting data filesystem dirs to original dirs"
    mount -o bind /mnt/primero-data/logs <%= node[:primero][:log_dir] %>
    mount -o bind /mnt/primero-data/couchdb /var/lib/couchdb
  else
    echo 'No primero-data filesystem found'
  fi
end script
