# cookbook/templates/default/primero_scheduler.service.erb
[Unit]
Description=Primero Task Scheduler
After=network.target

[Service]
Type=simple
Restart=always
PrivateTmp=yes
User=<%= node[:primero][:app_user] %>
Group=<%= node[:primero][:app_group] %>
WorkingDirectory=<%= node[:primero][:app_dir] %>
Environment="RAILS_ENV=<%= node[:primero][:rails_env] %>"
ExecStart=<%= node[:primero][:daemons_dir] %>/primero-scheduler-worker.sh

[Install]
WantedBy=multi-user.target
